# src/common/etc/profile.d/history.sh
#!/bin/sh

HISTTIMEFORMAT="%Y-%m-%d %T "
export HISTTIMEFORMAT
